# rootfs/etc/services.d/miner/run
#!/bin/sh

sleep 5
echo "* generating block..."

bitcoin-cli -conf=/etc/conf/node-1/bitcoin.conf generate 1
sleep 5
