# rootfs/etc/services.d/bitcoind-1/run
#!/bin/sh

echo "* Bitcoind node 1..."
exec /usr/bin/bitcoind \
    --printtoconsole \
    -datadir=/data/node-1 \
    -conf=/etc/conf/node-1/bitcoin.conf \
    -server \
    -walletnotify='/usr/bin/bitcoin-cli -conf=/etc/conf/node-1/bitcoin.conf gettransaction %s | tr "\n" " " | timeout --signal=SIGTERM 1 sed "s/$/\n/"  > /var/run/bitcoind/events-in' \
    -blocknotify='timeout --signal=SIGTERM 1 echo %s > /var/run/bitcoind/events-in'
